# iso-profile/airootfs/usr/lib/systemd/system/grub-btrfs.path
[Unit]
Description=Watch for snapshot changes

[Path]
PathModified=/etc/grub-btrfs/config
PathModified=/.snapshots

[Install]
WantedBy=multi-user.target