# iso-profile-arm64/airootfs/etc/systemd/system/systemd-networkd-wait-online.service.d/50-disable.conf
[Unit]
Description=Disabled for boot speed
ConditionPathExists=/dev/null
