# iso-profile/airootfs/etc/systemd/system/iso-auto-update.timer
[Unit]
Description=Weekly Auto Update

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=4h

[Install]
WantedBy=timers.target