# iso-profile/airootfs/etc/systemd/system/qdrant-grpc.socket
[Unit]
Description=Qdrant gRPC Socket
PartOf=qdrant.service

[Socket]
ListenStream=6334
Accept=yes

[Install]
WantedBy=sockets.target