# iso-profile/airootfs/etc/systemd/system/update-vm-nftables.timer
[Unit]
Description=Periodic VM Network Detection

[Timer]
OnBootSec=30
OnUnitActiveSec=5m
Persistent=true

[Install]
WantedBy=timers.target