# iso-profile/airootfs/etc/systemd/system/gdm.service.d/99-memory.conf
[Service]
MemoryHigh=192M
MemoryMax=256M
