# iso-profile-arm64/airootfs/etc/systemd/system/NetworkManager.service.d/99-memory.conf
[Service]
MemoryHigh=96M
MemoryMax=128M
