# iso-profile-arm64/airootfs/etc/systemd/system/iso-security-advisory.timer
[Unit]
Description=Weekly security advisory check

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
