# iso-profile-arm64/airootfs/etc/systemd/system/iso-btrfs-maintenance.service
[Unit]
Description=Btrfs maintenance — trim, dedup, balance
ConditionPathExists=/

[Service]
Type=oneshot
ExecStart=/usr/lib/iso/btrfs-maintenance.sh
StandardOutput=journal
StandardError=journal
