# iso-profile-arm64/airootfs/etc/systemd/system/iso-btrfs-maintenance.timer
[Unit]
Description=Weekly Btrfs maintenance

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=6h

[Install]
WantedBy=timers.target
