# iso-profile/airootfs/etc/systemd/system/gnome-shell.service.d/99-memory.conf
[Service]
MemoryHigh=384M
MemoryMax=512M
